# _legacy/rabbitmq/root/etc/services.d/rabbitmq/run
#!/usr/bin/execlineb -P
if { s6-svwait -t 5000 -u /var/run/s6/services/resolver }
with-contenv
s6-setuidgid rabbitmq
rabbitmq-server
